# dot_config/systemd/user/mbsync@.service
[Unit]
Description=Mailbox synchronization for %i

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync %i
ExecStartPost=/usr/bin/notmuch new
